# root/etc/services.d/syncthing/run
#!/usr/bin/with-contenv sh

exec syncthing -home=/srv/config -gui-address="0.0.0.0:$SYNCTHING_GUI_PORT" -no-browser -no-restart
